# build/profile/airootfs/etc/systemd/system/sway-start.service
[Unit]
Description=Start Sway Window Manager
After=graphical.target
Requires=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/sway
Restart=always
Environment=DISPLAY=:0
User=%u

[Install]
WantedBy=graphical.target
