# build/profile/airootfs/etc/systemd/system/wireless-scanner.service
[Unit]
Description=Wireless Network Scanner Service
DefaultDependencies=no
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wireless-scanner
ExecStartPre=/usr/bin/npm install
ExecStart=/usr/bin/node app.js
Restart=always
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target