# build/profile/airootfs/etc/systemd/system/wireless-scanner.service
[Unit]
Description=Wireless Network Scanner Service
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wireless-scanner
Environment=PATH=/usr/bin:/bin:/usr/local/bin
ExecStartPre=/usr/bin/npm install
ExecStart=/usr/bin/node app.js
Restart=always
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target