# build/profile/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
#ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin root - $TERM
ExecStart=-/sbin/agetty --noclear %I $TERM
Type=idle
