# build/profile/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin BROWSER_OS --noclear %I $TERM
Type=idle
